# benchmark.service
[Service]
ExecStart=/home/phil/bin/benchmark
LoadCredential=github_token
LoadCredential=webhook_secret
