# benchmark.service
[Unit]
Description=Benchmark Queue/Runner
After=network.target

[Service]
User=benchmarker
ExecStart=/usr/local/bin/benchmark serve 0.0.0.0:80
Environment=PATH=/home/benchmarker/micromamba/envs/asv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
SetCredentialEncrypted=webhook_secret:...
SetCredentialEncrypted=github_token:...

[Install]
WantedBy=default.target
